# borgmatic/files/default/borgmatic.timer.j2
[Unit]
Description=Run borgmatic backup

[Timer]
OnCalendar={{ borgmatic.service.timer }}
Persistent=true
RandomizedDelaySec={{ borgmatic.service.rand_delay }}

[Install]
WantedBy=timers.target
